# supervisord-haproxy.conf
[program:setup-haproxy]
command=/setup-haproxy.sh
autostart=true
auto_restart=false

[program:haproxy]
command=/start-haproxy.sh
numprocs=1
autostart=true
autorestart=true
