# supervisord-etcdctl.conf
[program:etcdctl]
command=/start-etcdctl.sh
numprocs=1
autostart=true
autorestart=true
